# com.user.firecorners.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.firecorners</string>
    <key>ProgramArguments</key>
    <array>
        <string>REPLACE_WITH_FULL_PATH_TO_SCRIPT/run_firecorners.sh</string>
        <!-- Add any command line arguments below if desired -->
        <!-- <string>--threshold=5</string> -->
        <!-- <string>--cooldown=3.0</string> -->
        <!-- <string>--dwell=0.5</string> -->
        <!-- <string>--no-test</string> -->
        <!-- <string>--config=/path/to/custom/config.json</string> -->
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>~/Library/Logs/FireCorners/firecorners.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/FireCorners/firecorners.error.log</string>
</dict>
</plist> 